# payara.conf
# upstart config script
description "Automatically upload contents of git repos to Amazon S3."
author      "Zach Kelling"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    exec /var/apps/payara/server.py >> /var/log/payara.log 2>&1
end script
